# rel/pkg/rpm/riak.service
[Unit]
Description=Riak KV Database

[Service]
User=riak
ExecStart=/usr/sbin/riak foreground
ExecStop=/usr/sbin/riak stop
Type=simple
PIDFile=/run/riak/riak.pid
EnvironmentFile=-/etc/default/riak
RuntimeDirectory=riak

[Install]
WantedBy=multi-user.target
