# package/smartos/manifest_smartos.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="riak">
  <service name="application/riak" type="service" version="1">
    <create_default_instance enabled="false" />
    <single_instance />
    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default" />
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local" />
    </dependency>
    <dependency name="riak-epmd" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/network/riak-epmd:default" />
    </dependency>
    <method_context working_directory="/tmp" project="riak">
      <method_credential user="riak" group="riak" />
      <method_environment>
        <envvar name="HOME" value="/opt/local/lib/riak" />
        <envvar name="LOGNAME" value="riak" />
        <envvar name="PATH" value="/usr/bin:/bin:/opt/local/bin:/opt/local/sbin" />
        <envvar name="LD_PRELOAD_32" value="/lib/libumem.so.1" />
        <envvar name="LD_PRELOAD_64" value="/lib/64/libumem.so.1:/opt/local/lib/amd64/libgcc_s.so.1"/>
        <envvar name="UMEM_OPTIONS" value="allocator=best" />
      </method_environment>
    </method_context>
    <exec_method type="method" name="start" exec="/opt/local/sbin/riak start" timeout_seconds="60" />
    <exec_method type="method" name="stop" exec="/opt/local/sbin/riak stop" timeout_seconds="60" />
    <exec_method type="method" name="restart" exec="/opt/local/sbin/riak restart" timeout_seconds="60" />
    <property_group name="startd" type="framework">
      <propval name="duration" type="astring" value="contract" />
      <propval name="ignore_error" type="astring" value="core,signal" />
    </property_group>
    <stability value="Evolving" />
    <template>
      <common_name>
        <loctext xml:lang="C">Riak Data Store</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
